# image/scripts/bin/crontab.cron
* * * * * sh /var/run/check_crontab.sh >/dev/stdout 2>&1